# .config/systemd/user/aria2c.service
[Unit]
Description=aria2c autostart

[Service]
Type=simple
ExecStart=/usr/bin/aria2c

[Install]
WantedBy=default.target
