# .config/systemd/user/pull-github-keys.timer
[Unit]
Description=Run pull-github-keys.service every day

[Timer]
OnBootSec=10min
OnUnitActiveSec=24h

[Install]
WantedBy=timers.target

