# .config/systemd/user/pull-github-keys.service
[Unit]
Description=Pull GitHub Keys
Wants=network-online.target
After=network-online.target

[Service]
ExecStart=%h/.local/bin/pull-github-keys.sh

